# systemd/tunnel-client-grabowski.service.d/70-operator-dependency.conf.example
[Unit]
Wants=grabowski-operator.service
After=grabowski-operator.service
