# systemd/tunnel-client-grabowski.service.d/70-operator-dependency.conf.example
[Unit]
# Follow operator stop/restart propagation without a hard failure coupling.
Wants=grabowski-operator.service
After=grabowski-operator.service
PartOf=grabowski-operator.service
