# systemd/tunnel-client-grabowski.service.d/80-restart-budget.conf.example
[Unit]
StartLimitIntervalSec=5min
StartLimitBurst=5
