# systemd/grabowski-privileged-broker@.service
[Unit]
Description=Grabowski privileged action broker request
Documentation=file:/usr/local/share/doc/grabowski/control-plane.md

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/libexec/grabowski-privileged-broker
StandardInput=socket
StandardOutput=socket
StandardError=journal
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=tmpfs
BindReadOnlyPaths=-/home/alex/.local/state/grabowski/recovery/last-server-recovery.json
BindReadOnlyPaths=-/home/alex/.local/state/grabowski/operator-kill-switch
BindReadOnlyPaths=-/home/alex/repos/.weltgewebe-worktrees
BindReadOnlyPaths=-/home/alex/worktrees
BindReadOnlyPaths=-/home/alex/repos/.semantah-worktrees
BindReadOnlyPaths=-/home/alex/repos/.heimlern-worktrees
BindReadOnlyPaths=-/home/alex/repos/.operator-redundancy-worktrees
BindReadOnlyPaths=-/home/alex/repos/.hauski-worktrees
BindReadOnlyPaths=-/home/alex/repos/.worktree-target-quarantine
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX
TimeoutStartSec=3700
