# systemd/grabowski-operator-watchdog.service.example
# Productive semantic watchdog: verify the live event loop and FastMCP session-
# creation lock without opening a competing session, plus one isolated stdio
# lifecycle from the exact deployed runtime. Before recovery restart it extracts
# one bounded sealed-memfd stack dump through an atomic private-file replace.
[Unit]
Description=Grabowski operator semantic watchdog
After=grabowski-operator.service
ConditionPathExists=%h/.local/libexec/grabowski/component_watchdog.py

[Service]
Type=oneshot
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 %h/.local/libexec/grabowski/component_watchdog.py --component operator --service grabowski-operator.service --runtime-root %h/.local/share/grabowski-mcp --module grabowski_operator --host 127.0.0.1 --port 18181 --failure-threshold 2 --max-restarts 3 --restart-window 900 --backoff-base 60 --backoff-max 900
# Exit 1 (failure observed / restart deferred) is routine evidence, not a unit
# failure. Exits 2 (indeterminate), 3 (budget exhausted) and 4 (restart did not
# recover) keep marking the unit failed.
SuccessExitStatus=1
# Worst case: probe + systemctl restart + bounded recovery re-probe loop.
TimeoutStartSec=90
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/grabowski
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
SystemCallArchitectures=native
