# systemd/grabowski-coding-agent-probe.service.example
[Unit]
Description=Refresh the advisory Grabowski coding-agent catalog
Documentation=file://%h/.local/share/grabowski/docs/coding-agent-routing-v2.md
ConditionPathExists=%h/bin/agent-route
ConditionPathExists=%h/.config/grabowski/coding-agent-probe-scheduler-router.sha256
ConditionPathExists=%h/.config/grabowski/coding-agent-catalog.json

[Service]
Type=oneshot
ExecStart=%h/.local/libexec/grabowski/coding_agent_probe_scheduler.py
TimeoutStartSec=3min
TimeoutStopSec=10s
KillMode=mixed
UMask=0077
Environment=PYTHONUNBUFFERED=1
Environment=PATH=%h/bin:%h/.local/bin:%h/.cargo/bin:%h/.npm-global/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=OPENAI_API_KEY=
Environment=ANTHROPIC_API_KEY=
Environment=XAI_API_KEY=
Environment=GEMINI_API_KEY=
Environment=GOOGLE_API_KEY=
Environment=OPENROUTER_API_KEY=
Environment=AZURE_OPENAI_API_KEY=
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/grabowski/coding-agent-router
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
