# systemd/grabowski-operator.service.example
[Unit]
Description=Grabowski MCP trusted-owner operator
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/.local/share/grabowski-mcp/.venv/bin/python -m grabowski_operator --transport streamable-http --host 127.0.0.1 --port 18181
Restart=on-failure
RestartSec=2
TimeoutStartSec=20
TimeoutStopSec=30
KillMode=mixed
UMask=0022
Environment=PYTHONUNBUFFERED=1
Environment=GRABOWSKI_TRUSTED_OWNER=1
NoNewPrivileges=no
PrivateTmp=no
ProtectSystem=off
ProtectHome=no
TasksMax=infinity
LimitNOFILE=1048576

[Install]
WantedBy=default.target
