# systemd/grabowski-coding-agent-probe.timer.example
[Unit]
Description=Refresh the advisory Grabowski coding-agent runtime state before it expires

[Timer]
OnBootSec=5min
OnUnitActiveSec=45min
RandomizedDelaySec=3min
AccuracySec=1min
Persistent=true
Unit=grabowski-coding-agent-probe.service

[Install]
WantedBy=timers.target
