# systemd/grabowski-runtime-retention.timer.example
[Unit]
Description=Periodically apply bounded Grabowski runtime retention

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=30s
Persistent=true
Unit=grabowski-runtime-retention.service

[Install]
WantedBy=timers.target
