# systemd/grabowski-watchdog.timer.example
[Unit]
Description=Run Grabowski semantic runtime watchdog

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=1s
RandomizedDelaySec=0
Persistent=false
Unit=grabowski-watchdog.service

[Install]
WantedBy=timers.target
