# systemd/grabowski-operator-watchdog.timer.example
[Unit]
Description=Run Grabowski operator watchdog

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=5s
RandomizedDelaySec=3s
Persistent=true
Unit=grabowski-operator-watchdog.service

[Install]
WantedBy=timers.target
