# systemd/grabowski-privileged-broker@.service.d/recovery-source.conf
[Service]
ProtectHome=tmpfs
BindReadOnlyPaths=
BindReadOnlyPaths=/home/alex/.local/state/grabowski/recovery/last-server-recovery.json
BindReadOnlyPaths=-/home/alex/.local/state/grabowski/operator-kill-switch
