# systemd/grabowski-privileged-broker@.service.d/recovery-source.conf
[Service]
ProtectHome=tmpfs
BindReadOnlyPaths=
BindReadOnlyPaths=/home/alex/.local/state/grabowski/recovery/last-server-recovery.json
BindReadOnlyPaths=-/home/alex/.local/state/grabowski/operator-kill-switch
BindReadOnlyPaths=-/home/alex/repos/.weltgewebe-worktrees
BindReadOnlyPaths=-/home/alex/worktrees
BindReadOnlyPaths=-/home/alex/repos/.semantah-worktrees
BindReadOnlyPaths=-/home/alex/repos/.heimlern-worktrees
BindReadOnlyPaths=-/home/alex/repos/.operator-redundancy-worktrees
BindReadOnlyPaths=-/home/alex/repos/.hauski-worktrees
BindReadOnlyPaths=-/home/alex/repos/.worktree-target-quarantine
